# systemd/aniremind-venv.service
[Unit]
Description=My Thing

[Service]
ExecStart=%h/.config/venvs/aniremind/bin/aniremind

[Install]
WantedBy=default.target
