# systemd/aniremind.service
[Unit]
Description=My Thing

[Service]
ExecStart=%h/.local/bin/aniremind

[Install]
WantedBy=default.target
